# macos/com.jamesdowzard.bose-control.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<!--
  LaunchAgent for the v2 menu-bar app. The resident process is a SILENT menu-bar
  agent (LSUIElement) that does NOT poll the headphones — it only touches the link
  on menu-open or a BT connect/disconnect event. KeepAlive simply keeps the menu-bar
  item present across logouts/crashes; it does not imply any background BT activity.
  Install via `macos/build.sh --install` — NOT installed automatically.
-->
<dict>
	<key>Label</key>
	<string>com.jamesdowzard.bose-control</string>
	<key>Program</key>
	<string>/Applications/Bose Control.app/Contents/MacOS/Bose Control</string>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/Users/jamesdowzard/Library/Logs/bose-control.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/jamesdowzard/Library/Logs/bose-control.log</string>
</dict>
</plist>
